# .local/share/backups/automatic-backup.service
[Service]
Type=oneshot
ExecStart=/home/zerinol/.local/share/backups/run.sh
